# interviews/deploy/interview-app.service
[Unit]
Description=AIBadges interviews app (B2B validation study)
After=network.target

[Service]
Type=simple
User=aibadges
WorkingDirectory=/opt/aibadges/interviews
Environment=PORT=4620
Environment=NODE_ENV=production
Environment=DATA_DIR=/opt/aibadges/interviews-data
# APP_USER / APP_PASS: set in /etc/aibadges/interviews.env (chmod 600)
EnvironmentFile=/etc/aibadges/interviews.env
# phase run injects NVIDIA_API_KEY (and the rest of the /global path) at launch;
# the key never lands on disk.
ExecStart=/usr/bin/env phase run --app-id 075a8ab8-78d4-4f75-9fdd-a94ba7d1712e --env Development --path /global -- bun run start
Restart=on-failure
RestartSec=3

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/aibadges/interviews-data
PrivateTmp=true

[Install]
WantedBy=multi-user.target
